# migrate/consultor-horizon.conf
; =============================================================================
; consultor-horizon.conf — Supervisor para Horizon (queues) y Scheduler
; Copiar a: /etc/supervisor/conf.d/consultor-horizon.conf
; Aplicar:  sudo supervisorctl reread && sudo supervisorctl update
; =============================================================================

; -----------------------------------------------------------------------------
; Laravel Horizon — gestiona todos los queue workers
; Reemplaza los queue:work individuales. El dashboard de Horizon está
; disponible en /horizon (protegerlo con HorizonServiceProvider::gate())
; -----------------------------------------------------------------------------
[program:consultor-horizon]
process_name=%(program_name)s
command=php8.4 /var/www/consultor/artisan horizon
directory=/var/www/consultor
user=www-data
autostart=true
autorestart=true

; Cuándo considerar que el proceso arrancó bien (segundos)
startsecs=5

; Intentos de restart antes de marcar como FATAL
startretries=3

; Horizon hace su propio graceful shutdown con SIGTERM
stopsignal=SIGTERM
stopwaitsecs=3600

; Matar hijos si el proceso padre muere
stopasgroup=true
killasgroup=true

redirect_stderr=true
stdout_logfile=/var/www/consultor/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

; Variables de entorno (Horizon las toma del .env, pero por si acaso)
environment=LARAVEL_SAIL=""

; -----------------------------------------------------------------------------
; Laravel Scheduler — ejecuta tareas programadas cada minuto
; Alternativa más limpia que el cron de www-data
; -----------------------------------------------------------------------------
[program:consultor-scheduler]
process_name=%(program_name)s
command=php8.4 /var/www/consultor/artisan schedule:work
directory=/var/www/consultor
user=www-data
autostart=true
autorestart=true
startsecs=3
startretries=3
stopsignal=SIGTERM
stopwaitsecs=60
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/www/consultor/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3

; -----------------------------------------------------------------------------
; Grupo — permite operar ambos juntos:
;   sudo supervisorctl restart consultor:*
;   sudo supervisorctl status consultor:*
; -----------------------------------------------------------------------------
[group:consultor]
programs=consultor-horizon,consultor-scheduler
priority=999
